# extra/usr/lib/systemd/user/codile.service
[Unit]
Description=Codile

[Service]
ExecStart=/usr/share/codile/app/daemon.sh

[Install]
WantedBy=default.target